# Tools/live-build/cgi/live-build-cgi.crontab
# /etc/crontab - system-wide crontab

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m  h dom mon dow user command
*/1  * *   *   *   root /etc/cron.hourly/live-webhelper > /dev/null 2>&1
